# var/spool/cron/crontab/pi
0-57/4 * * * * /home/pi/kion.sh
1-58/4 7-23 * * * /home/pi/aquestalkpi/AquestTalkPi -f '/home/pi/kion.txt' | aplay
